# systemd/systor-collector.service
[Unit]
Description=systor collector — system metrics poller + alerter
After=network-online.target
Wants=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=120

[Service]
Type=simple
ExecStart=/usr/bin/python3 -m systor serve collector
WorkingDirectory=/opt/systor
Restart=always
RestartSec=10
TimeoutStopSec=10
KillMode=control-group
StandardOutput=append:/var/log/systor/collector.log
StandardError=append:/var/log/systor/collector.log

# Load all SYSTOR_* env vars
EnvironmentFile=-/etc/systor/systor.env

[Install]
WantedBy=default.target
