# systemd/systor-web.service
[Unit]
Description=systor web — dashboard on port 6677
After=network-online.target systor-collector.service
Wants=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=120

[Service]
Type=simple
ExecStart=/usr/bin/python3 -m systor serve web
WorkingDirectory=/opt/systor
Restart=always
RestartSec=10
TimeoutStopSec=10
KillMode=control-group
StandardOutput=append:/var/log/systor/web.log
StandardError=append:/var/log/systor/web.log
EnvironmentFile=-/etc/systor/systor.env

[Install]
WantedBy=default.target
